# installer/service/com.checkpointvcs.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.checkpointvcs.daemon</string>

    <!--
      checkpoint-daemon is a portable Node.js runtime; it runs daemon-bundle.cjs.
      /usr/local/bin/checkpoint-daemon is a symlink into the app bundle
      (see installer/macos/scripts/postinstall); node resolves it to its real
      location so the bundle's native modules load from the daemon dir. The
      bundle path must be the real path inside the app bundle.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/checkpoint-daemon</string>
        <string>/Applications/Checkpoint.app/Contents/Resources/daemon/daemon-bundle.cjs</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/tmp/checkpoint-daemon.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/checkpoint-daemon.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
    </dict>
</dict>
</plist>
